# supervisor.tudo.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=supervisord.log ; main log file; default $CWD/supervisord.log
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
user=root
nodaemon=True

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ;

[program:gunicorn]
command=/bin/bash -c '/tudo-api/start_api.sh'
directory=/tudo-api/
stdout_logfile=/tudo-api/gunicorn.log
stderr_logfile=/tudo-api/gunicorn.log
autostart=true
autorestart=true
startsecs=5
killasgroup=true
stopasgroup=true
stopsignal=QUIT
priority=998

[program:celery-worker]
command=/bin/bash -c '/tudo-api/start_cw.sh'
directory=/tudo-api/
stdout_logfile=/tudo-api/celery-worker.log
stderr_logfile=/tudo-api/celery-worker.log
autostart=true
autorestart=true
startsecs=5
killasgroup=true
stopasgroup=true
stopsignal=QUIT
priority=998

[program:celery-beat]
command=/bin/bash -c '/tudo-api/start_cb.sh'
directory=/tudo-api/
stdout_logfile=/tudo-api/celery-beat.log
stderr_logfile=/tudo-api/celery-beat.log
autostart=true
autorestart=true
startsecs=5
killasgroup=true
stopasgroup=true
stopsignal=QUIT
priority=998
